# discord-bridge.service.example
[Unit]
Description=Discord Bridge (persistent Claude Agent SDK sessions)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USER/projects/discord-bridge
ExecStart=/usr/bin/node bridge.mjs
Restart=on-failure
RestartSec=10

# The SDK reads ~/.claude/* for CLAUDE.md, agents, skills, hooks, MCP servers.
# Make sure HOME is set so the right config gets picked up.
Environment=HOME=/home/YOUR_USER
Environment=PATH=/home/YOUR_USER/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

# Discord token + any overrides (ACCESS_JSON, STATE_DIR, etc.)
EnvironmentFile=/home/YOUR_USER/projects/discord-bridge/.env

[Install]
WantedBy=default.target
