# jules-discord-bridge.service.example
[Unit]
Description=Jules Discord Bridge (persistent Jules SDK sessions)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USER/projects/jules-discord-bridge
ExecStart=/usr/bin/node bridge.mjs
Restart=on-failure
RestartSec=10

Environment=HOME=/home/YOUR_USER
Environment=PATH=/home/YOUR_USER/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

# Discord token + any overrides (ACCESS_JSON, STATE_DIR, etc.)
EnvironmentFile=/home/YOUR_USER/projects/jules-discord-bridge/.env

[Install]
WantedBy=default.target
